# cebulobot.cron
SHELL=/bin/bash
PYTHONPATH="${PYTHONPATH}:/root/cebulobot/"

0 * * * * root python3 /root/cebulobot/bot/boy.py >> /var/log/cebulobot.log 2>&1